# userdata/systemd-units/etc-kubernetes.automount
[Unit]
Wants=btrfs-create-subvolumes.service
After=DATA.mount
After=btrfs-create-subvolumes.service
Before=kubelet.service bootkube.service

[Automount]
Where=/etc/kubernetes

[Install]
WantedBy=local-fs.target
WantedBy=kubelet.service bootkube.service
