# userdata/systemd-units/var-etcd.automount
[Unit]
Wants=btrfs-create-subvolumes.service
After=DATA.mount
After=btrfs-create-subvolumes.service
Before=kubelet.service

[Automount]
Where=/var/etcd

[Install]
WantedBy=local-fs.target
WantedBy=kubelet.service
