# userdata/systemd/DATA.automount
[Automount]
Where=/DATA

[Install]
WantedBy=local-fs.target
