# userdata/systemd/kubelet.path
[Unit]
Description=Watch for kubeconfig

[Path]
PathExists=/etc/kubeconfig

[Install]
WantedBy=multi-user.target
