# userdata/systemd/var-lib-kubelet.mount
[Unit]
Before=kubelet.service

[Mount]
What=/dev/disk/by-label/DATA
Where=/var/lib/kubelet
Type=btrfs
Options=compress=lzo,subvol=/kubelet
